# system/etc/cron.d/pigpiod
# https://www.raspberrypi.org/forums/viewtopic.php?f=32&t=103752
@reboot /usr/local/bin/pigpiod
